# home-server/udp-relay-home-for-relay.service
[Unit]
Description=Custom UDP relay (home side, for the relay VM's Bedrock path), bridges the relay's public Bedrock port to local Geyser
After=network.target frpc.service
Requires=frpc.service

[Service]
Type=simple
Environment=LISTEN_PORT=19135
Environment=RELAY_AUTH_TOKEN=REPLACE_ME
ExecStart=/usr/local/bin/udp_relay_home.ts
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
